# etc/hotplug.d/iface/99-preventleak
#!/bin/sh
# /etc/hotplug.d/iface/99-prevent-leak

# https://nordvpn.com/tutorials/openwrt/openvpn/

if [ "$ACTION" = ifup ] || [ "$ACTION" = ifdown ]; then
	/etc/openvpn/preventleak.sh
fi
