# nordvpn/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

rm -f /etc/openvpn/_openvpn.lock &
/etc/openvpn/_watchdog &

exit 0
